# rules.d/60-cr48-gps.rules
# udev rules for starting gps
KERNEL=="ttyUSB3", RUN+="/opt/Qualcomm/bin/start_gps.sh"
